# litestream/rootfs/etc/services.d/litestream/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Litestream SQLite Backup add-on"

# Read configuration
REPLICA_TYPE=$(bashio::config 'replica_type')
BUCKET=$(bashio::config 'bucket')
PATH_PREFIX=$(bashio::config 'path')
SYNC_INTERVAL=$(bashio::config 'sync_interval')
RETENTION=$(bashio::config 'retention')
RETENTION_CHECK_INTERVAL=$(bashio::config 'retention_check_interval')

# Database path
DB_PATH="/config/home-assistant_v2.db"

# Check if database exists
if [ ! -f "$DB_PATH" ]; then
    bashio::log.error "Home Assistant database not found at $DB_PATH"
    bashio::log.error "Make sure Home Assistant has started at least once"
    sleep 30
    exit 1
fi

bashio::log.info "Replica type: ${REPLICA_TYPE}"
bashio::log.info "Bucket: ${BUCKET}"
bashio::log.info "Path: ${PATH_PREFIX}"

# Build litestream configuration
CONFIG_FILE="/tmp/litestream.yml"

if [ "$REPLICA_TYPE" = "gcs" ]; then
    GCS_CREDS=$(bashio::config 'gcs_credentials_json')

    if [ -z "$GCS_CREDS" ]; then
        bashio::log.error "GCS credentials JSON is required for GCS replica type"
        sleep 30
        exit 1
    fi

    # Write credentials to file
    printf '%s' "$GCS_CREDS" > /tmp/gcs-credentials.json
    export GOOGLE_APPLICATION_CREDENTIALS="/tmp/gcs-credentials.json"

    # Debug: verify credentials file
    bashio::log.info "Credentials file size: $(wc -c < /tmp/gcs-credentials.json) bytes"
    bashio::log.info "Credentials project: $(grep -o '"project_id"[^,]*' /tmp/gcs-credentials.json | head -1)"

    cat > "$CONFIG_FILE" << EOF
dbs:
  - path: ${DB_PATH}
    replicas:
      - type: gcs
        bucket: ${BUCKET}
        path: ${PATH_PREFIX}
        sync-interval: ${SYNC_INTERVAL}
        retention: ${RETENTION}
        retention-check-interval: ${RETENTION_CHECK_INTERVAL}
EOF

elif [ "$REPLICA_TYPE" = "s3" ]; then
    S3_ENDPOINT=$(bashio::config 's3_endpoint')
    S3_ACCESS_KEY=$(bashio::config 's3_access_key_id')
    S3_SECRET_KEY=$(bashio::config 's3_secret_access_key')
    S3_REGION=$(bashio::config 's3_region')

    if [ -z "$S3_ACCESS_KEY" ] || [ -z "$S3_SECRET_KEY" ]; then
        bashio::log.error "S3 access key and secret key are required for S3 replica type"
        sleep 30
        exit 1
    fi

    export AWS_ACCESS_KEY_ID="$S3_ACCESS_KEY"
    export AWS_SECRET_ACCESS_KEY="$S3_SECRET_KEY"

    cat > "$CONFIG_FILE" << EOF
dbs:
  - path: ${DB_PATH}
    replicas:
      - type: s3
        bucket: ${BUCKET}
        path: ${PATH_PREFIX}
        sync-interval: ${SYNC_INTERVAL}
        retention: ${RETENTION}
        retention-check-interval: ${RETENTION_CHECK_INTERVAL}
        region: ${S3_REGION}
EOF

    if [ -n "$S3_ENDPOINT" ]; then
        echo "        endpoint: ${S3_ENDPOINT}" >> "$CONFIG_FILE"
    fi
fi

bashio::log.info "Starting Litestream replication..."
exec /usr/local/bin/litestream replicate -config "$CONFIG_FILE"
